# variants/base_Xfce_buster_amd64/system-config/includes.chroot/etc/cron.d/night-mode
# automatic night mode: from blue to red color shift to protect sleep quality
* * * * *   root	/usr/local/bin/night-mode
