# features/package-list_system_tools/livefs-include/etc/systemd/system/cpupower-gui.service.d/override.conf
[Unit]
ConditionVirtualization=false

[Service]
SuccessExitStatus=255
