# variants/FSFW-Uni-Stick_KDE_stretch_amd64/config/includes.chroot/etc/systemd/system/tmp.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=basic.target

[Mount]
Where=/tmp
What=tmpfs
Type=tmpfs
Options=auto,nosuid,nodev

[Install]
WantedBy=live-config.service
