# features/config_desktop_cats/user-config/.config/autostart/desktop_cats.desktop
[Desktop Entry]
Type=Application
Name=Desktop Cats
Comment=have some cute cats chase the mouse cursor
Exec=/usr/local/bin/desktop-cats.sh
Terminal=false
NoDisplay=true
