# features/total_conversion_chicago95/livefs-include/usr/local/share/Chicago95/etc/skel/.config/autostart/welcome3.avi.desktop
[Desktop Entry]
Type=Application
Name=welcome3.avi
Comment=Windows 95 extra
Exec=sh -c "sleep 5; mpv --fullscreen /usr/local/share/Chicago95/welcome3.avi"
