# variants/FSFW-Uni-Stick_Xfce_buster_amd64/system-config/includes.chroot/etc/profile.d/stty-flow.sh
#!/bin/sh
# turn off stupid flow control and shut up.
stty -ixon -ixoff 2>/dev/null
