# variants/FSFW-Uni-Stick_Xfce_buster_amd64/system-config/includes.chroot/etc/profile.d/user-local-bin-path.sh
#!/bin/sh
case "$PATH" in
  */$HOME* ) ;;
  * ) PATH="$HOME/.local/bin:$PATH";
esac
