# features/config_terminal_pipx/livefs-include/etc/profile.d/pipx-global-environment.sh
export PIPX_HOME=/usr/local/pipx
export PIPX_BIN_DIR=/usr/local/bin
