# features/config_jupyter_notebook_service/user-config/.config/systemd/user/jupyter-notebook.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/jupyter-notebook --no-browser --NotebookApp.token=

