# features/config_dresden_schulproxy/livefs-include/etc/systemd/system/NetworkManager.service.d/override.conf
[Service]
ExecStopPost=-systemctl mask redsocks | systemd-cat -t dresden-schulen-proxy
