# variants/base_Xfce_buster_amd64/system-config/includes.chroot/etc/xdg/autostart/lightsonplus.desktop
[Desktop Entry]
Type=Application
Name=lightsonplus
Comment=disable the screensaver when playing videos
TryExec=/usr/local/libexec/lightsonplus/lightson+
Exec=/usr/local/libexec/lightsonplus/lightson+
Terminal=false
NoDisplay=true
