# features/config_apt_unattended-upgrades/livefs-include/etc/systemd/system/apt-daily-upgrade.timer.d/override.conf
# Expand schedule for applying package updates
[Timer]
OnBootSec=10min
RandomizedDelaySec=0
Persistent=false
