# variants/base_Xfce_buster_amd64/user-config/.config/autostart/Terminal.desktop
[Desktop Entry]
Type=Application
Name=Terminal
Comment=Xfce-Terminal in drop-down mode
Exec=xfce4-terminal --drop-down
