# features/package-list_dictionaries/livefs-include/etc/systemd/system/dictd.service.d/override.conf
[Service]
# OOM Killing the dictionary server not a problem
OOMScoreAdjust=300
