# features/config_automatic_night_mode/user-config/.config/systemd/user/night-mode.timer
[Unit]
Description=automatic night mode (timer)

[Timer]
OnCalendar=*-*-* *:*:42

[Install]
WantedBy=default.target
