# features/config_automatic_night_mode/user-config/.config/systemd/user/night-mode.service
[Unit]
Description=automatic night mode: from blue to red color shift to protect sleep quality
Requires=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/night-mode
