# features/config_xfce_light-dark-switcher/user-config/.config/autostart/light-dark-switcher.desktop
[Desktop Entry]
Name=Light/Dark-Switcher
Type=Application
Exec=light-dark-switcher
StartupNotify=false
