# features/config_apt_unattended-upgrades/livefs-include/etc/systemd/system/apt-daily-upgrade.service.d/override.conf
[Service]
# Do not run auto-upgrades when booting without persistent storage
# because system might run out of RAM with larger updates.
ExecCondition=sh -c '[ -d /run/live/persistence/*/ ]'
