# variants/base_Xfce_buster_amd64/system-config/includes.chroot/etc/systemd/system/console-setup.service.d/override.conf
# workaround for https://bugs.debian.org/cgi-bin/bugreport.cgi?bug=846256
[Unit]
After=systemd-tmpfiles-setup.service
