# features/config_terminal_environment/livefs-include/etc/profile.d/games-path.sh
#!/bin/sh
case "$PATH" in
  */games* ) ;;
  * ) PATH="$PATH:/usr/games" ;;
esac
