# features/config_terminal_environment/livefs-include/etc/profile.d/systemd-no-pager-in-tmux.sh
# No systemd pager in tmux
[ -n "$TMUX" ] && export SYSTEMD_PAGER=
