# features/service_prevent_out_of_memory_freezes_nohang/livefs-include/etc/systemd/user/dbus.service.d/override.conf
[Service]
# make it a bit less likely whole session gets shot down
# (can't go lower than 100 unfortunately)
OOMScoreAdjust=100
