# features/config_apt_unattended-upgrades/livefs-include/etc/systemd/system/apt-daily.timer.d/override.conf
# Expand schedule for refreshing package lists
[Timer]
OnBootSec=5min
RandomizedDelaySec=2min
