# features/config_compressed_ram/livefs-include/etc/systemd/system/systemd-zram-setup@zram1.service.d/override.conf
[Service]
ExecStartPost=/bin/sh -xc 'd=$(mktemp -d); mount "$1" "$d"; chmod -c 1777 "$d"; umount "$d"; rmdir "$d"' _ /dev/%i
