# variants/FSFW-Uni-Stick_Xfce_stretch_amd64/system-config/includes.chroot/etc/profile.d/mellow-beep.sh
[ $TERM = "linux" ] && setterm --bfreq 230 --blength 42
