# hotfixes/OSEG-Makerstick_3057k_Xfce_bookworm-amd64/persistence-data/linux-userdata/user/.config/autostart/light-dark-switcher.desktop
[Desktop Entry]
Name=Theme switcher
Type=Application
Exec=light-dark-switcher
