# features/config_tuptime_logging/livefs-include/etc/systemd/system/tuptime-sync.service.d/override.conf
# Less log pollution
[Service]
LogLevelMax=warning
