# features/config_devilspie_window_automator/user-config/.config/autostart/devilspie_window_automater.desktop
[Desktop Entry]
Type=Application
Name=Devilspie
Comment=find windows and perform actions on them
Exec=devilspie
Terminal=false
NoDisplay=true
